# DockerExample/crontab
* * * * * sh /root/entry.sh >/proc/self/fd/1 2>/var/error.log
