# DockerExample/crontab
* * * * * /root/runner.sh >/proc/self/fd/1 2>/var/error.log
